# examples/OpenStack/dope.service
[Unit]
Description=FixNet service

[Service]
ExecStart=/opt/net.sh

[Install]
WantedBy=multi-user.target
